# systemd/framegrabber.timer
[Unit]
Description=Poll for Steam Frame availability every 2 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
# Run a missed poll on resume (e.g. after sleep) instead of waiting for the next slot.
Persistent=true

[Install]
WantedBy=timers.target
